# deb/DEBIAN/prerm
#!/bin/bash


update-rc.d -f fcgid remove
